# deploy_conf/supervisor.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=/tmp/supervisor_wx_robot_cmd.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)
chown=dlnu:dlnu        ; socket file uid:gid owner

[supervisord]
logfile=/tmp/supervisor_wx_robot_cmd.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisor_wx_robot_cmd.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor_wx_robot_cmd.sock ; use a unix:// URL  for a unix socket


[program:wx_robot_cmd]
user=dlnu
directory=/home/dlnu/WX-ROBOT
environment=PYTHONPATH=/home/dlnu/WX-ROBOT
command=/home/dlnu/Envs/wx_robot/bin/python3 weixin.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/dlnu/WX-ROBOT/logs/supervisor/out_log.log
stdout_logfile_maxbytes=500MB
stdout_logfile_backups=50
stderr_logfile=/home/dlnu/WX-ROBOT/logs/supervisor/err_log.log
exitcodes=0
loglevel=info
